# services/nginx/supervisor.conf
[program:nginx]
command=nginx -c baselayer/conf/nginx.conf -p . -g "daemon off;"
