# services/nginx/supervisor.conf
[program:nginx]
command=nginx -c baselayer/services/nginx/nginx.conf -p . -g "daemon off;"
stdout_logfile=log/nginx.log
redirect_stderr=true
