# services/nginx/supervisor.conf
[program:nginx]
command=nginx -c baselayer/services/nginx/nginx.conf -p . -g "daemon off;"
# First nginx up early
priority=50
